# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/mulmul/bin/uwsgi --ini /srv/mulmul/.config/uwsgi/uwsgi-app.ini
